# generated_bundle/systemd/hermes-user-test.service
[Unit]
Description=Hermes Agent for user_test
After=network.target

[Service]
Type=simple
User=hmx_user_test
Group=hmx_user_test
WorkingDirectory=/home/hmx_user_test
Environment=HOME=/home/hmx_user_test
Environment=HERMES_HOME=/home/hmx_user_test/.hermes
ExecStart=/usr/local/bin/hermes gateway run --replace
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
